# docker/conf/supervisor/apache.conf
[program:apache]
command=apache2-foreground

